# debian/prerm
#!/bin/sh
set -e

#DEBHELPER#

SERVICE_NAME="web-rust-print"

case "$1" in
    remove|upgrade|deconfigure)
        # 停止服务但不删除配置
        systemctl stop "$SERVICE_NAME" 2>/dev/null || true
        systemctl disable "$SERVICE_NAME" 2>/dev/null || true
        ;;
esac

exit 0
