# debian/web-rust-print.service
[Unit]
Description=Web Rust Print Server (Axum + SQLite)
After=network.target

[Service]
Type=simple
User=root

ExecStart=/usr/bin/web_rust_print

WorkingDirectory=/var/lib/web-rust-print

# 启动超时时间
TimeoutStartSec=30

# 崩溃/退出后自动重启
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
